# background.conf
## service_script.conf

[supervisord]  
nodaemon=false

[program:worker_script] 
command=python3 Worker/scheduling.py  
autorestart=true
stderr_logfile=/dev/stdout 
stderr_logfile_maxbytes = 0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes = 0
## same setting for 2nd service
# [program:django_service] 
# command=gunicorn --bind 0.0.0.0:8000 fampay_api.wsgi
# autostart=true
# autorestart=true
# stderr_logfile=/dev/stdout
# stderr_logfile_maxbytes = 0
# stdout_logfile=/dev/stdout
# stdout_logfile_maxbytes = 0